# webservice/tei-annotator.service
# Systemd service unit for the TEI Annotator webservice.
#
# Setup:
#   1. Adjust the three placeholder values below:
#        USER          — the OS user to run as (e.g. the user who cloned the repo)
#        WORKING_DIR   — absolute path to the webservice/ directory
#        UV_PATH       — output of: which uv
#      Quick one-liner (run from the repo root, logged in as the service user):
#        sudo sed -i \
#          -e "s|USER|$(whoami)|g" \
#          -e "s|WORKING_DIR|$(pwd)/webservice|g" \
#          -e "s|UV_PATH|$(which uv)|g" \
#          /etc/systemd/system/tei-annotator.service
#   2. Copy to the systemd directory:
#        sudo cp webservice/tei-annotator.service /etc/systemd/system/
#   3. Enable and start:
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now tei-annotator
#   4. Check status and follow logs:
#        sudo systemctl status tei-annotator
#        journalctl -u tei-annotator -f

[Unit]
Description=TEI Annotator webservice
After=network.target

[Service]
Type=simple
User=USER
WorkingDirectory=WORKING_DIR
EnvironmentFile=WORKING_DIR/.env
# Remove any stale PID file before each start so the built-in
# process-replacement logic does not interfere with systemd restarts.
# The leading '-' makes this step non-fatal if the file doesn't exist
# or the path is not yet accessible.
ExecStartPre=-rm -f WORKING_DIR/.server.pid
ExecStart=UV_PATH run python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
